# pimd.service.in
[Unit]
Description=PIM-SM/SSM multicast routing daemon
Documentation=man:pimd
Documentation=man:pimctl
Documentation=man:pimd.conf
# ConditionPathExists=@SYSCONFDIR@/pimd.conf
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
EnvironmentFile=-@SYSCONFDIR@/default/pimd
ExecStart=@SBINDIR@/pimd -ns $PIMD_OPTS

[Install]
WantedBy=multi-user.target
