# debian/pimd.upstart
description "pimd - PIM-SM v2 dynamic multicast routing daemon"

start on filesystem
stop on runlevel [016]

expect fork
respawn

pre-start exec test -x /usr/sbin/pimd || { stop; exit 0; }
exec /usr/sbin/pimd

# Daemon can not reload, only restart
reload signal 0
